# debian/postrm
#!/bin/sh
set -e

#DEBHELPER#

if [ "$1" = "purge" ] ; then
	update-rc.d mkenv remove >/dev/null
fi
